# functional/fixtures/units/conflict.2.service
[Unit]
Description=Test Unit

[Service]
ExecStart=/usr/bin/sleep 120

[X-Fleet]
X-Conflicts=conflict.*.service
