# functional/fixtures/units/global.service
[Service]
ExecStart=/bin/bash -c "while true; do echo Hello, Global World!; sleep 1; done"

[X-Fleet]
Global=true
