# functional/fixtures/units/goodbye.service
[Service]
ExecStart=/bin/bash -c "while true; do echo Goodbye, World!; sleep 1; done"
