# functional/fixtures/units/pong.2.service
[Unit]
Description=PONG

[Service]
ExecStart=/bin/bash -c "while true; do echo pong; sleep 1; done"

[X-Fleet]
X-ConditionMachineOf=ping.2.service
