# functional/fixtures/units/single@.service
[Service]
ExecStart=/bin/bash -c "while true; do echo 'Hello, World (@%i)!'; sleep 1; done"
