# examples/pong.service
[Unit]
Description=PONG

[Service]
ExecStart=/bin/bash -c "while true; do echo \"pong\"; sleep 1; done"

[X-Fleet]
MachineOf=ping.service
