# examples/socket-activated.socket
[Unit]
Description=Socket-Activated Web Service

[Socket]
ListenStream=8080
