# examples/web.service
[Unit]
Description=Web Service

[Service]
ExecStart=/usr/bin/python -m SimpleHTTPServer
